# ansible/roles/node_exporter/templates/local.d/node_exporter.stop
#!/bin/bash

start-stop-daemon --user {{ node_exporter_monitoring_user }} --pidfile /run/node_exporter.pid --stop

