# ansible/roles/node_exporter/templates/local.d/node_exporter.start
#!/bin/bash

start-stop-daemon --user {{ node_exporter_monitoring_user }} --pidfile /run/node_exporter.pid --background --make-pidfile --start --exec /usr/local/bin/node_exporter -- --collector.textfile.directory={{ node_exporter_textfile_directory }}

