# ansible/roles/monit/templates/local.d/monit.start
#!/bin/bash

# {{ ansible_managed }}

rm -f {{ monit_pid_file }}
{{ monit_exec_file_path }}

