# ansible/roles/monit/templates/local.d/monit.stop
#!/bin/bash

# {{ ansible_managed }}

{{ monit_exec_file_path }} quit
sleep 3
killall -s KILL {{ monit_exec_file_path }}
