# ansible/roles/monit_collector/templates/local.d/monit_collector.start
#!/bin/bash

su {{ monit_collector_monitoring_user }} -c "{{ monit_collector_script_name }} --daemon"

