# exporters_install_scripts/packer_example/root_sync_preinstall/etc/cron.d/cron_test
* * * * * monitoring date > /tmp/monit_cron_test
