# exporters_install_scripts/packer_example/root_sync_preinstall/etc/cron.d/monit_not_running
* * * * * monitoring /usr/local/bin/monit_not_running.sh > /dev/null 2>&1
