# exporters_install_scripts/packer_example/root_sync_preinstall/etc/cron.d/monit_config_update
@reboot root /usr/local/scripts/monit/monit_config_update.sh > /dev/null 2>&1

